# documentation/connectivity/eod-capture@XEUR.timer
# eod-capture@XEUR.timer — daily Eurex close-capture trigger (WS 1G, ADR 0032).
#
# Fires the eod-capture@XEUR.service shortly after the Eurex regular close. The trigger time is
# a FIXED local time stated in the EXCHANGE timezone (Europe/Berlin), never server-local — Eurex
# closes at 17:30 CET/CEST, so 18:15 Europe/Berlin is a safe upper bound on the regular close.
#
# The fixed OnCalendar time is only the *trigger*: the runner resolves the EXACT close instant
# from the 1J exchange-calendar (session_close), so a half-day early close or a holiday is
# handled by the resolver (a holiday is a clean no-op), NOT by editing this timer. The timezone
# is explicit so the trigger tracks the Eurex close across DST, not the server's clock.
#
# Persistent=true: a fire missed while the box was down runs on next boot; the catch-up fire
# reconstructs the gap day through the run-state ledger (idempotent by construction).
#
# Install: copy to ~/.config/systemd/user/ and `systemctl --user enable --now eod-capture@XEUR.timer`.

[Unit]
Description=Daily Eurex (XEUR) close-capture timer

[Timer]
# 18:15 Europe/Berlin, Mon–Fri — a safe upper bound on the 17:30 Eurex close. The runner skips
# any day the XEUR calendar marks a non-session, so weekends/holidays are no-ops, not misfires.
OnCalendar=Mon..Fri 18:15 Europe/Berlin
Persistent=true
# A small randomized delay so multiple calendars firing near each other don't thundering-herd.
RandomizedDelaySec=120
Unit=eod-capture@XEUR.service

[Install]
WantedBy=timers.target
